# pi-setup/mk3-settings-watcher.service
[Unit]
Description=MK3 Settings Command Executor
After=mixxx.service

[Service]
Type=simple
User=pi
Environment=DISPLAY=:99
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/python3 /home/pi/mixx-mk3/pi-setup/mk3-settings-watcher.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
