# pi-setup/mk3-mouse-daemon.service
[Unit]
Description=MK3 Mouse Mode Daemon
After=xvfb.service mixxx.service
Wants=xvfb.service

[Service]
Type=simple
User=pi
Environment=DISPLAY=:99
ExecStart=/usr/bin/python3 /home/pi/mixx-mk3/pi-setup/mk3-mouse-daemon.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
