# pi-setup/mk3-screen-daemon.service
[Unit]
Description=MK3 Screen Mirroring Daemon
# Started by udev when MK3 is plugged in, stopped when removed
StopWhenUnneeded=yes

[Service]
Type=simple
ExecStart=/usr/local/bin/mk3-screen-daemon --resolution 960x544 --fps 15 --no-partial
Restart=on-failure
RestartSec=2
